# backend/deploy/secretvoip-mail.service
[Unit]
Description=SecretVoIP Mail backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=svm
Group=svm
WorkingDirectory=/opt/secretvoip-mail/backend
EnvironmentFile=/opt/secretvoip-mail/backend/.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secretvoip-mail

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes

[Install]
WantedBy=multi-user.target
